# extras/service-files/opensrf-websockets.service.j2
[Unit]
Description=OpenSRF Websockets

[Service]
Type=simple
User=opensrf
Group=opensrf
ExecStart=/usr/local/bin/opensrf-websockets --max-clients {{websockets_max_clients}}
