# extras/service-files/opensrf-router.service.j2
[Unit]
Description=OpenSRF Router

[Service]
Type=simple
User=opensrf
Group=opensrf
ExecStart=/openils/bin/opensrf-router
